# kube_cluster/roles/master/templates/bridge.service.j2
[Unit]
Description=sysctl minus p

[Service]
ExecStart=/sbin/sysctl -p

[Install]
WantedBy=multi-user.target
